# templates/default/systemd/consul.service.erb
# -*- mode: conf-unix -*-
[Unit]
Description=Consul Agent
Wants=basic.target
After=basic.target network.target

[Service]
EnvironmentFile=/etc/conf.d/<%= @instance %>
User=<%= @run_user %>
Group=<%= @run_group %>
ExecStart=<%= Chef::Consul.active_binary(node) %> agent -config-dir <%= @config_dir %> <%= @extra_options %>
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
