# templates/default/consul.conf.erb
description "Consul Service Discovery Platform"

emits consul-up

start on runlevel [2345]
stop on runlevel [!2345]

# set max open files
limit nofile <%= @soft_limit %> <%= @hard_limit %>

script
  if [ -f <%= node['consul']['etc_config_dir'] %> ]; then
    . <%= node['consul']['etc_config_dir'] %>
  fi
  export GOMAXPROCS=${GOMAXPROCS}
  CMD="<%= Chef::Consul.active_binary(node) %> agent -config-dir <%= node['consul']['config_dir'] %>"
  LOGFILE="<%= @consul_logfile %>"
  exec $CMD >> "$LOGFILE"
end script

post-start exec initctl emit consul-up

kill signal INT

respawn
respawn limit 10 10
